# deploy/pdp/services/curio-poststart.service
[Unit]
Description=Curio Post-Start Tasks
After=curio-ready.service
Wants=curio-ready.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/curio-poststart.sh

[Install]
WantedBy=multi-user.target
